# crontab
* * * * * python3 /app/emailNotificator.py >> /var/log/out.log 2>&1
* * * * * python3 /app/telegramNotificator.py >> /var/log/out.log 2>&1